# rootfs/rootfs_tree/etc/inittab
::sysinit:/init
::askfirst:/bin/sh
::ctrlaltdel:/sbin/reboot
